# ansible/templates/backup.service.j2
[Unit]
Description=Restic backup for restaurants app

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup-restaurants.sh
